# SystemdFiles/radiko-recorder.timer
[Unit]
Description=radiko-recorder

[Timer]
Unit=radiko-recorder.service
OnCalendar=*-*-* 18:00:00

[Install]
WantedBy=timers.target
